# packaging/deb/DEBIAN/prerm
#!/bin/sh
set -e

if command -v systemctl >/dev/null 2>&1; then
  systemctl stop mcpanel.service || true
  systemctl disable mcpanel.service || true
  systemctl daemon-reload || true
fi

exit 0
